# ovirt-cockpit-sso.service
[Unit]
Description=oVirt-Cockpit SSO service

[Service]
Restart=always
ExecStartPre=/usr/sbin/remotectl certificate --ensure --user=ovirt --group=ovirt --selinux-type=etc_t
ExecStart=/usr/share/ovirt-cockpit-sso/start.sh
User=ovirt
Group=ovirt

[Install]
WantedBy=multi-user.target
